# worker-monitor.service
[Unit]
Description=Worker Activity Monitor Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/your/worker-monitor/main.py
Restart=always
User=your_username
WorkingDirectory=/path/to/your/worker-monitor

[Install]
WantedBy=multi-user.target
